# admin-dashboard/control-center-admin.service
[Unit]
Description=TCC Custom - Admin Dashboard
Documentation=https://github.com/rick0490/tcc-custom
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/tcc-custom/admin-dashboard
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
EnvironmentFile=/root/tcc-custom/admin-dashboard/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
