# MagicMirror-flyer/magic-mirror-flyer.service
[Unit]
Description=Magic Mirror Flyer Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/tournament-dashboard/MagicMirror-flyer
ExecStart=/usr/bin/node --run server
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
