# tournament-signup/tournament-signup.service
[Unit]
Description=Tournament Signup Web App
Documentation=https://github.com/yourusername/tournament-dashboard
After=network.target magic-mirror-match.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/tournament-dashboard/tournament-signup
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening (optional)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
